# twingate-indicator.service
[Unit]
Description=Twingate Connection Indicator
Documentation=https://www.twingate.com/docs/linux
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=%h/.local/bin/twingate-indicator daemon
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Run in background
KillMode=process

[Install]
WantedBy=default.target
