# twingate-tray.service
[Unit]
Description=Twingate System Tray Indicator
Documentation=https://www.twingate.com/docs/linux
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/local/bin/twingate-tray
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Run in background
KillMode=process

[Install]
WantedBy=default.target
